# config/supervisord.conf
#Set supervisord log location
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

##################################################################
# The bellow enables supervisorctl in the container, starts here
; rpc interface for supervisorctl

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/run/supervisord.sock

[supervisorctl]
serverurl=unix:///run/supervisord.sock

# The Above enables supervisorctl in the container, ends here
##################################################################

# Set api to run with nodemon
[program:app]
command=nodemon server.js
directory=/app/
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=1
 
